# 4.1 bike back up code vcu/system/log-compress.service
[Unit] 
Description=Compress Specific Log File
[Service] 
Type=oneshot 
ExecStart=/home/root/script_files/compress-logs.sh
